# train_mix3d.conf
[supervisord]

[program:Mix3D]
command = python mix3d general.gpus=1
startsecs = 0
autorestart = false
